# systemd/docker-gfs-plugin.socket
[Unit]
Description=Docker Volume Driver for gfs volumes

[Socket]
ListenStream=/run/docker/plugins/gfs.sock

[Install]
WantedBy=sockets.target
